# ironclaw-deployment/systemd/ironclaw.service
[Unit]
Description=IronClaw Secure Telegram Agent
Documentation=https://github.com/nearai/ironclaw
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/ironclaw

# Main process
ExecStart=/home/pi/ironclaw/target/release/ironclaw daemon
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=60

# Environment
Environment=RUST_LOG=info
Environment=IRONCLAW_CONFIG=/home/pi/.ironclaw/settings.toml

# =============================================================================
# SECURITY HARDENING
# =============================================================================

# Prevent privilege escalation
NoNewPrivileges=true

# Filesystem restrictions
ProtectSystem=strict
ProtectHome=read-only

# Allow write access only to specific directories
ReadWritePaths=/home/pi/ironclaw-notes
ReadWritePaths=/var/log/ironclaw

# Protect kernel and hardware
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true

# Restrict device access
PrivateDevices=true

# Restrict network namespaces (but allow network access for Telegram API)
RestrictNamespaces=~user ~pid ~net ~uts ~mnt ~cgroup ~ipc

# Memory protections
MemoryDenyWriteExecute=true

# Restrict realtime scheduling
RestrictRealtime=true

# Restrict address families to only what's needed
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# System call filter - allow only needed syscalls
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Capability restrictions - drop all capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# Temporary directory
PrivateTmp=true

# Hide other users' processes
ProtectProc=invisible
ProcSubset=pid

# Limit resource usage
MemoryMax=512M
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ironclaw

[Install]
WantedBy=multi-user.target
